# app/templates/app/config/mina/templates/systemd.service.erb
[Unit]
Description=<%= app_sys_id %>
After=network.target

[Service]
ExecStart=<%= current_release %>/server.js
Restart=always
User=<%= user %>
Group=<%= group %>

Environment=PATH=/usr/bin:/usr/local/bin
Environment=NODE_ENV=production
<% scheme = enable_ssl ? 'https://' : 'http://' %>
Environment=API_HOST=<%= scheme + public_domain %>
Environment=APP_PORT=<%= app_port %>
Environment=API_PORT=<%= api_port %>

WorkingDirectory=<%= current_release %>

[Install]
WantedBy=multi-user.target
